# systemd/tampering-check@.service.example
[Unit]
Description=Tampering Check Service for %I
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/tampering-check.sh %I
Restart=on-failure
StandardOutput=journal
StandardError=journal
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
RuntimeDirectory=tampering-check
RuntimeDirectoryMode=0700
StateDirectory=tampering-check
StateDirectoryMode=0700
ReadWritePaths=/var/spool/postfix/maildrop

[Install]
WantedBy=multi-user.target

