# systemd/tampering-check@.service.example
[Unit]
Description=Tampering Check Service for %I
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/tampering-check.sh %I
Restart=on-failure
StandardOutput=journal
StandardError=journal
ProtectSystem=strict
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/log/tampering-check /var/lock /var/spool/postfix/maildrop

[Install]
WantedBy=multi-user.target

